# docker/rc.local
#! /bin/bash

cd /root/ && sudo ./configure_collectd_${DEPLOYERL_MODE}.sh
